# docker/multi_container/supervisord.conf
[supervisord]
nodaemon=true

[program:python]
user=root
directory=/home/oee_webapp
autostart = true
autorestart = true
stopwaitsecs=600
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=gunicorn -b :8000 -m 007 "app:create_app()"

[program:celeryworker]
user=root
directory=/home/oee_webapp
autostart = true
autorestart = true
stopwaitsecs=600
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=celery -A celery_worker.celery_app worker

[program:celerybeat]
user=root
directory=/home/oee_webapp
autostart = true
autorestart = true
stopwaitsecs=600
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command= celery -A celery_worker.celery_app beat